# mailcatcher/supervisord.conf
[supervisord]
nodaemon=true

[program:mailcatcher]
command=/usr/bin/mailcatcher --foreground --ip 0.0.0.0
redirect_stderr=true
